# resources/etc/supervisord.d/httpd.conf
[program:httpd]
command=/usr/sbin/httpd -D "FOREGROUND" -k start
redirect_stderr=true
