# resources/etc/supervisord.d/php.conf
[program:php]
command=/usr/sbin/php-fpm
autostart=true
